# deploy/ansible/roles/docker/files/prodenv/omnileads-prodenv.service
[Unit]
Description=Omnileads devenv Application Service
Requires=docker.service
After=docker.service

[Service]
User={{ docker_user }}
Type=oneshot
RemainAfterExit=yes
Environment=COMPOSE_FILE={{ deploy_location }}/prodenv_stack.yml
ExecReload=/usr/local/bin/docker-compose up -d
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
